# deployment/systemd/cloudflared.service
[Unit]
Description=Cloudflare Tunnel Client
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-apps/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as root for tunnel (or create dedicated cloudflared user)
User=root
Group=root

# Configuration file location
Environment=TUNNEL_CONFIG=/etc/cloudflared/config.yml

# Start the tunnel
ExecStart=/usr/local/bin/cloudflared tunnel --config /etc/cloudflared/config.yml run

# Restart configuration - tunnel should always be running
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadOnlyPaths=/etc/cloudflared

# Resource limits
MemoryMax=128M
CPUQuota=30%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

[Install]
WantedBy=multi-user.target
